# htest-server/source/config/system/htest-server.service
[Unit]
Description=Hardware tester web server
Requires=network.target
Requires=mysql.service
After=network.target
After=mysql.service

[Service]
Type=simple
Restart=always
ExecStart=/usr/bin/node /usr/share/htest-server/index.js

[Install]
WantedBy=multi-user.target
